# crontab
0 * * * * python /home/nick/weather/slurp.py
